# deploy/supervisor/zycdn-dns-sync.conf
# ============================================================================
# ShieldFlow DNS 同步服务 Supervisor 配置 (shieldflow-dns-sync.conf)
# 管理: shieldflow-dns-sync
# 部署路径: /etc/supervisord.d/shieldflow-dns-sync.conf
# 进程组名: shieldflow-dns-sync
# ============================================================================
[group:shieldflow-dns-sync]
programs=shieldflow-dns-sync
priority=10

[program:shieldflow-dns-sync]
command=/usr/local/shieldflow/bin/dns-sync -config /etc/shieldflow/dns-sync.yaml
directory=/usr/local/shieldflow
user=root
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=15
stopsignal=TERM
stdout_logfile=/var/log/shieldflow/dns-sync.stdout.log
stderr_logfile=/var/log/shieldflow/dns-sync.stderr.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5
priority=10
