# deploy/supervisor/zycdn-log-server.conf
# ============================================================================
# ShieldFlow 日志服务器 Supervisor 配置 (shieldflow-log-server.conf)
# 管理: shieldflow-log-server
# 部署路径: /etc/supervisord.d/shieldflow-log-server.conf
# 进程组名: shieldflow-log-server
# ============================================================================
[group:shieldflow-log-server]
programs=shieldflow-log-server
priority=10

[program:shieldflow-log-server]
command=/usr/local/shieldflow/bin/log-server -config /etc/shieldflow/log-server.yaml
directory=/usr/local/shieldflow
user=root
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=15
stopsignal=TERM
stdout_logfile=/var/log/shieldflow/log-server.stdout.log
stderr_logfile=/var/log/shieldflow/log-server.stderr.log
stdout_logfile_maxbytes=200MB
stdout_logfile_backups=20
stderr_logfile_maxbytes=200MB
stderr_logfile_backups=20
priority=10
