# deploy/supervisor/zycdn-master.conf
# ============================================================================
# ShieldFlow 主控节点 Supervisor 配置 (shieldflow-master.conf)
# 管理: shieldflow-backend (API) + shieldflow-grpc-server (gRPC)
# 部署路径: /etc/supervisord.d/shieldflow-master.conf
# 进程组名: shieldflow-master
# 命令: supervisorctl reread && supervisorctl update
# ============================================================================
[group:shieldflow-master]
programs=shieldflow-backend,shieldflow-grpc-server
priority=10

[program:shieldflow-backend]
command=/usr/local/shieldflow/bin/backend -config /etc/shieldflow/backend.yaml
directory=/usr/local/shieldflow
user=root
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=15
stopsignal=TERM
stdout_logfile=/var/log/shieldflow/backend.stdout.log
stderr_logfile=/var/log/shieldflow/backend.stderr.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=10
environment=GIN_MODE="release"
priority=10

[program:shieldflow-grpc-server]
command=/usr/local/shieldflow/bin/grpc-server -config /etc/shieldflow/grpc.yaml
directory=/usr/local/shieldflow
user=root
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=15
stopsignal=TERM
stdout_logfile=/var/log/shieldflow/grpc-server.stdout.log
stderr_logfile=/var/log/shieldflow/grpc-server.stderr.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=10
priority=20
