# deploy/supervisor/zycdn-edge.conf
# ============================================================================
# ShieldFlow 边缘节点 Supervisor 配置 (shieldflow-edge.conf)
# 管理: shieldflow-edge
# 部署路径: /etc/supervisord.d/shieldflow-edge.conf
# 进程组名: shieldflow-edge
# ============================================================================
[group:shieldflow-edge]
programs=shieldflow-edge
priority=10

[program:shieldflow-edge]
command=/usr/local/shieldflow/bin/edge -config /etc/shieldflow/edge.yaml
directory=/usr/local/shieldflow
user=root
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=15
stopsignal=TERM
stdout_logfile=/var/log/shieldflow/edge.stdout.log
stderr_logfile=/var/log/shieldflow/edge.stderr.log
stdout_logfile_maxbytes=200MB
stdout_logfile_backups=20
stderr_logfile_maxbytes=200MB
stderr_logfile_backups=20
environment=GIN_MODE="release"
priority=10
