# deploy/supervisor-dns-sync.conf
[program:shieldflow-dns-sync]
command=/opt/shieldflow/bin/dns-sync -config /etc/shieldflow/dns-sync.yaml
directory=/opt/shieldflow
user=shieldflow
group=shieldflow
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=15
stdout_logfile=/var/log/shieldflow/dns-sync.stdout.log
stderr_logfile=/var/log/shieldflow/dns-sync.stderr.log
environment=GIN_MODE="release"
